# A-Stock-Analysis/com.huigu.astock-analysis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.huigu.astock-analysis</string>
  <key>ProgramArguments</key><array>
    <string>/usr/bin/python3</string>
    <string>{{SCRIPT_DIR}}/run_daily.py</string>
  </array>
  <key>WorkingDirectory</key><string>{{SCRIPT_DIR}}</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>12</integer>
      <key>Minute</key><integer>10</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>17</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>{{REPO_DIR}}/.git</string>
  </array>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>{{SCRIPT_DIR}}/logs/launchd.log</string>
  <key>StandardErrorPath</key><string>{{SCRIPT_DIR}}/logs/launchd.err.log</string>
</dict>
</plist>
